# systemd/daily-update.service
[Unit]
Description=Daglig systemuppdatering
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/dnf upgrade -y
